# deploy/nv.service
[Unit]
Description=NV Master Agent Harness
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=doppler run --project nova --config prd -- node %h/.local/lib/nova-ts/packages/daemon/dist/index.js
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30
Environment=NODE_ENV=production
Environment=HOME=%h
Environment=PATH=%h/.local/share/pnpm:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=%h/.local/lib/nova-ts
LimitNOFILE=4096
MemoryMax=2G
TasksMax=512

[Install]
WantedBy=default.target
