# deploy/nv.service
[Unit]
Description=NV Master Agent Harness
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=all
ExecStart=%h/.local/bin/doppler run --fallback=true --project nova --config prd -- %h/.local/bin/nv-daemon
Restart=on-failure
RestartSec=5s
WatchdogSec=60
TimeoutStopSec=30
Environment=RUST_LOG=info
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
LimitNOFILE=4096
MemoryMax=2G
TasksMax=128

[Install]
WantedBy=default.target
