# deploy/nv.service
[Unit]
Description=NV Master Agent Harness
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=%h/.local/bin/nv-daemon
Restart=on-failure
RestartSec=5s
WatchdogSec=60
TimeoutStopSec=30
EnvironmentFile=-%h/.nv/env
Environment=RUST_LOG=info
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
LimitNOFILE=4096
MemoryMax=512M
TasksMax=128

[Install]
WantedBy=default.target
