# relays/teams/nv-teams-relay.service
[Unit]
Description=Nova Teams Webhook Relay
After=network-online.target nv.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/doppler run --fallback=true --project nova --config prd -- /usr/bin/python3 %h/nv/relays/teams/server.py
Restart=on-failure
RestartSec=10s
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
